# disable8021xAutomaticConnection/org.povertyactionlab.disable8021xAutomaticConnection.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.povertyactionlab.disable8021xAutomaticConnection</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/disable8021xAutomaticConnection.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
